# system/systemd.bak/paccache-clean.service
[Unit]
Description=Clean Pacman Cache

[Service]
ExecStart=/usr/bin/paccache -r

